# scripts/birdlog.service
[Unit]
Description=BirdLog Server
After=network.target

[Service]
Type=simple
User=bullton
WorkingDirectory=/home/bullton/bird
Environment=PORT=3005
Environment=BIRDLOG_HOST=0.0.0.0
Environment=NODE_ENV=production
ExecStart=/usr/bin/node server/dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/bullton/bird/logs/server.log
StandardError=append:/home/bullton/bird/logs/server.err.log

[Install]
WantedBy=multi-user.target